# dotfiles/config/autostart/screensaver-settings.desktop
[Desktop Entry]
Name=Screensaver
Comment=Set screensaver timeouts
Exec=xset s 600 dpms 900 0 0
